# debian/extra/systemd-udevd.service.d/hc-utils.conf
[Service]
IPAddressAllow=169.254.169.254/32
